# modules/launch/host/generate-id-block/mkosi.extra/usr/local/lib/systemd/system/generate-id-block.service
[Unit]
Description=Generate id-block.b64 and id-auth.b64 for guest launch
DefaultDependencies=no
After=calculate-measurement.service
Wants=calculate-measurement.service
Before=launch-guest.service

[Service]
Type=oneshot
# ID block metadata — override any of these to change what is embedded in the
# signed block without modifying the script.  Defaults match the benchmark
# platform configuration documented in DESIGN.md.
Environment=ID_BLOCK_FAMILY_ID=sev-certify-fam0
Environment=ID_BLOCK_IMAGE_ID=sev-certify-img0
Environment=ID_BLOCK_GUEST_SVN=48
Environment=ID_BLOCK_POLICY=0xb0000
ExecStart=/usr/local/lib/scripts/generate_id_block.py
StandardOutput=journal+console
StandardError=journal+console
LogExtraFields="SEV_VERSION=3.0.0-0" "SNPHOST_TEST=3.0.0-0"
