# modules/launch/host/generate-id-block/mkosi.extra/usr/local/lib/systemd/system/generate-id-block.service
[Unit]
Description=Patch and sign id-block.b64, generating id-auth.b64
DefaultDependencies=no
After=calculate-measurement.service
Wants=calculate-measurement.service
Before=launch-guest.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/scripts/generate_id_block.py
StandardOutput=journal+console
StandardError=journal+console
LogExtraFields="SEV_VERSION=3.0.0-0" "SNPHOST_TEST=3.0.0-0"
